# pkgs/applications/window-managers/swhkd/swhkd.service
[Unit]
Description=swhkd hotkey daemon
BindsTo=default.target

[Service]
Type=simple
ExecStart=@out@

[Install]
WantedBy=default.target
